# root/etc/s6-overlay/s6-rc.d/deluged-set-folders/up
/etc/s6-overlay/scripts/deluged-set-folders
