# root/etc/s6-overlay/s6-rc.d/deluged-initial-setup/up
/etc/s6-overlay/scripts/deluged-initial-setup
